# scripts/robot-ble.service
[Unit]
Description=Robot BLE Control Bridge (Nordic UART Service -> Arduino over USB serial)
After=bluetooth.target
Wants=bluetooth.target

[Service]
Type=simple
# Adjust the path if you cloned the repo somewhere other than /opt/robot.
ExecStart=/usr/bin/python3 /opt/robot/scripts/ble_server.py --serial /dev/ttyACM0
Restart=always
RestartSec=5
User=root
Group=root

StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
